# etc/profile.d/ld_lib.sh
export LD_LIBRARY_PATH=$EPREFIX/usr/lib

